# templates/rc.local.erb
#!/bin/sh -e
# File Managed by Puppet
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.
#
# With Puppet is managed by the content of <%= scope.lookupvar('rclocal::params::configdir') %>

for command in $(ls -v1 <%= scope.lookupvar('rclocal::params::configdir') %> ) ; do
    <%= scope.lookupvar('rclocal::params::configdir') %>/$command
done

<% if @osfamily == 'RedHat' %>
touch /var/lock/subsys/local
<% end %>
<% if @osfamily == 'Debian' %>
exit 0
<% end %>
